# img_builder/packages/popusetnet_config/popusetnet.service
[Unit]
Description=popusetNET interface
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash /usr/local/bin/popusetnet.sh

[Install]
WantedBy=multi-user.target